# GRC-HAT/07-APPLICATION-LAYER/03-templates/wazuh/ossec.conf
<!-- ossec.conf — Hardened Wazuh Manager Configuration
     NIST: SI-7 (integrity), AU-2 (event logging), SI-4 (monitoring), IR-4 (incident handling)
     Deploy to: /var/ossec/etc/ossec.conf (manager or agent)
     Restart: systemctl restart wazuh-manager
-->
<ossec_config>

  <!-- ─── Manager / Agent Global Settings ─────────────────────────────────── -->
  <global>
    <!-- WHY: JSON alerts integrate with SIEM (Splunk/Sentinel) — structured data -->
    <jsonout_output>yes</jsonout_output>
    <!-- WHY: Log all alerts to file as backup — AU-11 retention requirement -->
    <alerts_log>yes</alerts_log>
    <!-- WHY: Log format includes hostname for multi-host correlation -->
    <logall>no</logall>
    <!-- Alert threshold: level 7+ fires to alerts.log; level 3+ to syslog -->
    <log_alert_level>3</log_alert_level>
    <!-- WHY: Email alerts for critical events (configure mail server below) -->
    <email_notification>no</email_notification>
    <!-- Uncomment and configure for email alerts:
    <email_notification>yes</email_notification>
    <smtp_server>smtp.company.com</smtp_server>
    <email_from>wazuh@company.com</email_from>
    <email_to>soc@company.com</email_to>
    <email_alert_level>12</email_alert_level>
    -->
  </global>

  <!-- ─── Remote (Agent Communication) ────────────────────────────────────── -->
  <remote>
    <!-- WHY: UDP 1514 is standard Wazuh agent-manager communication -->
    <connection>secure</connection>
    <port>1514</port>
    <protocol>tcp</protocol>
    <!-- WHY: TCP is more reliable than UDP for high-volume environments -->
    <queue_size>131072</queue_size>
  </remote>

  <!-- ─── Syscheck (File Integrity Monitoring) ─────────────────────────────── -->
  <!-- NIST SI-7: Detect unauthorized changes to software and information -->
  <syscheck>
    <!-- WHY: Disabled by default in some installs — SI-7 requires FIM -->
    <disabled>no</disabled>

    <!-- Scan frequency: 43200 = every 12 hours for scheduled paths -->
    <frequency>43200</frequency>

    <!-- WHY: Alert when new files appear in monitored directories -->
    <alert_new_files>yes</alert_new_files>

    <!-- WHY: Don't auto-ignore frequently changing files — investigate them -->
    <auto_ignore frequency="10" timeframe="3600">no</auto_ignore>

    <!-- Realtime monitoring: critical paths that must alert immediately
         WHY: Changes to these paths are pre-attack signals or active compromise -->
    <!-- T1003 credential dump target — realtime alert required -->
    <directories realtime="yes" report_changes="yes" check_all="yes">/etc/passwd</directories>
    <!-- T1003 credential dump target — realtime alert required -->
    <directories realtime="yes" report_changes="yes" check_all="yes">/etc/shadow</directories>
    <!-- T1548 privilege escalation path — any change is an alert -->
    <directories realtime="yes" report_changes="yes" check_all="yes">/etc/sudoers</directories>
    <directories realtime="yes" report_changes="yes" check_all="yes">/etc/sudoers.d</directories>
    <!-- T1098 SSH backdoor via config change -->
    <directories realtime="yes" report_changes="yes" check_all="yes">/etc/ssh/sshd_config</directories>

    <!-- Scheduled monitoring: important but less time-critical -->
    <!-- T1053 scheduled task persistence -->
    <directories check_all="yes">/etc/crontab</directories>
    <directories check_all="yes">/etc/cron.d</directories>
    <directories check_all="yes">/var/spool/cron</directories>
    <!-- T1565 DNS manipulation -->
    <directories check_all="yes">/etc/hosts</directories>
    <directories check_all="yes">/etc/resolv.conf</directories>
    <!-- T1610 Kubernetes control plane config integrity -->
    <directories check_all="yes">/etc/kubernetes</directories>
    <!-- SI-7 Wazuh self-protection — detect tampering with security tool -->
    <directories check_all="yes">/var/ossec/etc</directories>
    <!-- T1574 binary hijacking detection -->
    <directories check_all="yes">/usr/bin</directories>
    <directories check_all="yes">/usr/sbin</directories>
    <directories check_all="yes">/bin</directories>
    <directories check_all="yes">/sbin</directories>

    <!-- Ignore high-churn directories that generate noise without signal -->
    <ignore>/proc</ignore>
    <ignore>/sys</ignore>
    <ignore>/dev</ignore>
    <ignore>/run</ignore>
    <ignore>/tmp</ignore>
    <ignore type="sregex">\.log$</ignore>
    <ignore type="sregex">\.tmp$</ignore>
    <ignore type="sregex">\.pid$</ignore>
    <ignore type="sregex">\.swp$</ignore>
    <!-- Ignore Wazuh's own log files -->
    <ignore>/var/ossec/logs</ignore>
    <ignore>/var/ossec/queue</ignore>
    <ignore>/var/ossec/stats</ignore>

    <!-- Registry monitoring (Windows agents) -->
    <!-- WHY: Registry run keys are primary Windows persistence mechanism — T1547 -->
    <windows_registry>HKEY_LOCAL_MACHINE\Software\Microsoft\Windows\CurrentVersion\Run</windows_registry>
    <windows_registry>HKEY_LOCAL_MACHINE\Software\Microsoft\Windows\CurrentVersion\RunOnce</windows_registry>
    <windows_registry>HKEY_LOCAL_MACHINE\System\CurrentControlSet\Services</windows_registry>
  </syscheck>

  <!-- ─── Rootcheck ──────────────────────────────────────────────────────── -->
  <!-- NIST SI-7: Detect rootkits and hidden processes -->
  <rootcheck>
    <!-- WHY: Rootcheck detects hidden processes, ports, files — T1014 rootkit -->
    <disabled>no</disabled>
    <check_files>yes</check_files>
    <check_trojans>yes</check_trojans>
    <check_dev>yes</check_dev>
    <check_sys>yes</check_sys>
    <check_pids>yes</check_pids>
    <check_ports>yes</check_ports>
    <check_if>yes</check_if>
    <!-- Scan every 24 hours — rootkits don't come and go quickly -->
    <frequency>86400</frequency>
    <rootkit_files>/var/ossec/etc/shared/rootkit_files.txt</rootkit_files>
    <rootkit_trojans>/var/ossec/etc/shared/rootkit_trojans.txt</rootkit_trojans>
    <system_audit>/var/ossec/etc/shared/system_audit_rcl.txt</system_audit>
    <system_audit>/var/ossec/etc/shared/system_audit_ssh.txt</system_audit>
  </rootcheck>

  <!-- ─── Log Analysis ──────────────────────────────────────────────────── -->
  <!-- NIST AU-2: Event logging for all required audit events -->
  <localfile>
    <!-- Linux authentication events — SSH, sudo, PAM -->
    <log_format>syslog</log_format>
    <location>/var/log/auth.log</location>
  </localfile>

  <localfile>
    <log_format>syslog</log_format>
    <location>/var/log/syslog</location>
  </localfile>

  <localfile>
    <!-- Linux audit daemon — syscall level events -->
    <log_format>audit</log_format>
    <location>/var/log/audit/audit.log</location>
  </localfile>

  <localfile>
    <!-- Kubernetes audit log (control plane nodes) -->
    <log_format>json</log_format>
    <location>/var/log/kubernetes/audit/audit.log</location>
    <label key="source">k8s-audit</label>
  </localfile>

  <localfile>
    <!-- Docker daemon events -->
    <log_format>syslog</log_format>
    <location>/var/log/dpkg.log</location>
  </localfile>

  <!-- Windows Event Log (for Windows agents) -->
  <localfile>
    <log_format>eventlog</log_format>
    <location>Security</location>
  </localfile>

  <localfile>
    <log_format>eventlog</log_format>
    <location>System</location>
  </localfile>

  <!-- ─── Active Response ───────────────────────────────────────────────── -->
  <!-- NIST IR-4: Incident handling — automated response to detected threats -->
  <active-response>
    <!-- WHY: Block IPs with 10+ failed logins — stops brute force at host level -->
    <!-- Response: firewall-drop (iptables block for 600 seconds = 10 minutes) -->
    <command>firewall-drop</command>
    <location>local</location>
    <rules_group>authentication_failures</rules_group>
    <timeout>600</timeout>
    <repeated_offenders>30,60,120</repeated_offenders>
  </active-response>

  <active-response>
    <!-- WHY: Alert SOC on rootcheck detections — these do not auto-block -->
    <command>host-deny</command>
    <location>local</location>
    <rules_id>510,511,512</rules_id>
    <timeout>0</timeout>
  </active-response>

  <!-- ─── Cluster (Manager) ─────────────────────────────────────────────── -->
  <!-- Uncomment for multi-node Wazuh cluster deployment
  <cluster>
    <name>wazuh</name>
    <node_name>master-01</node_name>
    <node_type>master</node_type>
    <key>CLUSTER_KEY_HERE</key>
    <port>1516</port>
    <bind_addr>0.0.0.0</bind_addr>
    <nodes>
      <node>master-01</node>
    </nodes>
    <hidden>no</hidden>
    <disabled>no</disabled>
  </cluster>
  -->

</ossec_config>
