# OSI-MODEL-SECURITY/07-APPLICATION-LAYER/03-templates/splunk/inputs.conf
# inputs.conf — Gold-standard Splunk inputs for security monitoring
# NIST: AU-2 (event logging), AU-12 (audit generation), SI-4 (monitoring)
# Deploy to: $SPLUNK_HOME/etc/apps/search/local/inputs.conf
# Restart Splunk after changes: systemctl restart splunk

# ─── HTTP Event Collector (HEC) ───────────────────────────────────────────────
# WHY: HEC is the preferred ingest path for structured JSON from modern sources.
# Each source type gets its own token for access control and index routing.
# Never share tokens across source types — you lose the ability to revoke per-source.

[http]
disabled = 0
# NIST AU-12: Force encrypted transport — never send auth tokens over plaintext
enableSSL = 1
port = 8088
# Restrict to specific IPs in production — open 0.0.0.0 only in lab
# dedicatedIoThreads = 2

# Token: Application JSON logs (microservices, APIs)
[http://jsa-app-logs]
disabled = 0
index = main
sourcetype = _json
# WHY: App logs contain auth events, errors, business transactions — AU-2 core
token = ${HEC_TOKEN_APP}

# Token: Suricata network IDS alerts
[http://jsa-suricata]
disabled = 0
index = security
sourcetype = suricata
# WHY: Suricata EVE JSON provides L3/L4/L7 network threat detection — SI-4
token = ${HEC_TOKEN_SURICATA}

# Token: Wazuh HIDS/FIM alerts
[http://jsa-wazuh]
disabled = 0
index = security
sourcetype = wazuh
# WHY: Wazuh provides host-level FIM, rootcheck, active response events — SI-7
token = ${HEC_TOKEN_WAZUH}

# Token: Kubernetes container stdout/stderr
[http://jsa-k8s-containers]
disabled = 0
index = k8s
sourcetype = kube:container:logs
# WHY: Container logs capture application behavior and runtime errors — AU-2
token = ${HEC_TOKEN_K8S_CONTAINERS}

# Token: Kubernetes audit log (API server)
[http://jsa-k8s-audit]
disabled = 0
index = k8s
sourcetype = kube:apiserver:audit
# WHY: K8s audit log is the source of truth for cluster access — AU-12, AC-6
token = ${HEC_TOKEN_K8S_AUDIT}

# ─── File Monitor Inputs (Linux — Universal Forwarder) ───────────────────────
# WHY: File monitors capture logs from systems that don't support HEC natively.
# Required on Linux endpoints where syslog forwarding is not configured.

# Linux auth.log — SSH logins, sudo events, PAM events
[monitor:///var/log/auth.log]
disabled = 0
index = linux_secure
sourcetype = linux_secure
# WHY: auth.log contains every authentication event on Linux — AU-2 baseline
# NIST: AC-7 (unsuccessful logon attempts), AU-2 (event logging)
crcSalt = <SOURCE>
followTail = 0

# Linux syslog — general system events
[monitor:///var/log/syslog]
disabled = 0
index = main
sourcetype = syslog
# WHY: syslog captures system-level events beyond auth — service starts/stops, etc.
crcSalt = <SOURCE>

# Linux audit daemon — auditd (syscall-level audit)
[monitor:///var/log/audit/audit.log]
disabled = 0
index = linux_secure
sourcetype = linux:audit
# WHY: auditd provides syscall-level visibility — privileged operations, file access
# NIST: AU-12 (audit generation) — generates events at OS level
crcSalt = <SOURCE>

# Wazuh JSON alerts (direct file monitor — alternative to HEC)
[monitor:///var/ossec/logs/alerts/alerts.json]
disabled = 0
index = security
sourcetype = wazuh
# WHY: Direct file monitor as fallback when HEC not configured
crcSalt = <SOURCE>

# Kubernetes audit log (on control plane node)
[monitor:///var/log/kubernetes/audit/audit.log]
disabled = 0
index = k8s
sourcetype = kube:apiserver:audit
# WHY: K8s audit log — every API request, who made it, what they did
# NIST: AU-12 (audit generation), AC-6 (least privilege verification)
crcSalt = <SOURCE>

# ─── Windows Inputs (Windows Universal Forwarder) ────────────────────────────
# WHY: Windows Security Event Log is the primary audit source on Windows endpoints.
# EventIDs below map to NIST AU-2 required audit events.

[WinEventLog://Security]
disabled = 0
index = wineventlog
whitelist = 4624,4625,4627,4648,4663,4670,4672,4673,4698,4720,4722,4724,4728,4732,4756,4768,4769
# 4624: Successful logon           — AU-2 authentication events
# 4625: Failed logon               — AC-7 unsuccessful logon attempts
# 4648: Logon with explicit creds  — T1078 lateral movement signal
# 4663: Object access              — AC-3 access enforcement audit
# 4672: Special privileges         — AC-6 least privilege monitoring
# 4698: Scheduled task created     — T1053 persistence detection
# 4720: User account created       — AC-2 account management
# 4728: Member added to global grp — AC-2 group membership changes
# 4768/4769: Kerberos TGT/TGS     — T1558 Kerberoasting detection

[WinEventLog://System]
disabled = 0
index = wineventlog
whitelist = 7045,7036
# 7045: New service installed      — T1543 persistence / T1050
# 7036: Service state changed      — operational audit

[WinEventLog://Application]
disabled = 0
index = wineventlog

# PowerShell Script Block Logging (requires GPO to enable)
[WinEventLog://Microsoft-Windows-PowerShell/Operational]
disabled = 0
index = wineventlog
whitelist = 4103,4104
# 4103: Module logging             — command invocations
# 4104: Script block logging       — full PowerShell content — T1059.001
